# sv/rssbot/run
#!/usr/bin/env bash
cd ~/tr-rss-haber-izleme
exec .venv/bin/python src/bot.py
